# supervisord/sentry.conf
[program:sentry]
command=/var/praekelt/sentry/ve/bin/sentry --config=/var/praekelt/sentry/sentry_conf.py start
environment=PATH="/var/praekelt/sentry/ve/bin"
autorestart=true
#TODO: Add logging
