# supervisord/sentry.conf
[program:sentry]
command=sentry start --config=/var/praekelt/sentry/sentry.conf.py
directory=/var/praekelt/sentry
environment=PATH="/var/praekelt/sentry/ve/bin"
autorestart=true
