# supervisord/sentry.conf
[program:sentry]
command=/var/praekelt/sentry/ve/bin/sentry start --config=/var/praekelt/sentry/sentry.conf.py
environment=PATH="/var/praekelt/sentry/ve/bin"
autorestart=true
#TODO: Add logging
